# config/supervisor/boss888-worker.conf
; ==========================================================================
; Boss888 / LevelUp OS — Supervisor Queue Worker Config
; ==========================================================================
;
; Copy to: /etc/supervisor/conf.d/boss888-worker.conf
; Then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start boss888-worker:*
;

[program:boss888-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/boss888/artisan queue:work redis --queue=tasks-high,tasks,tasks-low,default --sleep=3 --tries=4 --backoff=8,16,32,64 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/boss888/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=3600
