# config/deploy/levelup-worker.conf
; ════════════════════════════════════════════════════════════════════
; LevelUp Growth Platform — Supervisor Queue Worker Config
; Target: Digital Ocean Droplet (Ubuntu 24.04)
;
; Install:
;   sudo cp levelup-worker.conf /etc/supervisor/conf.d/levelup-worker.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start levelup-worker:*
;
; Check status:
;   sudo supervisorctl status levelup-worker:*
;   sudo supervisorctl tail levelup-worker_00 stdout
; ════════════════════════════════════════════════════════════════════

[program:levelup-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/levelup/artisan queue:work redis \
    --queue=tasks-high,tasks,tasks-low,default \
    --sleep=3 \
    --tries=4 \
    --backoff=8,16,32,64 \
    --max-time=3600 \
    --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/levelup/storage/logs/worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/www/levelup/storage/logs/worker-error.log
stderr_logfile_maxbytes=10MB
stopwaitsecs=3600
priority=999

; ── Laravel Scheduler (single process, runs artisan schedule:run every minute) ──

[program:levelup-scheduler]
command=bash -c 'while true; do php /var/www/levelup/artisan schedule:run >> /var/www/levelup/storage/logs/scheduler.log 2>&1; sleep 60; done'
process_name=%(program_name)s
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/levelup/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stopasgroup=true
killasgroup=true
priority=998
